# etc/com.chatbotchat.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Reference launchd agent for the chatbotchat daemon.

  You normally do NOT edit or copy this by hand: `cbc install-daemon` resolves the
  daemon's absolute path, writes the equivalent agent into ~/Library/LaunchAgents/,
  and loads it. This file documents what that agent looks like and serves as a
  manual fallback.

  Manual install: replace CHANGE_ME with the absolute path to chatbotchat-server
  (`which chatbotchat-server`), then:

    cp etc/com.chatbotchat.server.plist ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/com.chatbotchat.server.plist

  Logs land in ~/Library/Logs/chatbotchat.log (+ .err.log). launchd does not
  rotate them; add a newsyslog.d rule if they grow (see the README).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chatbotchat.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/CHANGE_ME/.cargo/bin/chatbotchat-server</string>
        <string>--port</string>
        <string>8484</string>
    </array>

    <!-- Start at login and restart on crash -> always-on. -->
    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <!-- Long-running, low-priority background service. -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Grace window before launchd escalates SIGTERM -> SIGKILL on stop/reboot,
         so the daemon can finish in-flight work rather than being hard-killed. -->
    <key>ExitTimeOut</key>
    <integer>25</integer>

    <key>StandardOutPath</key>
    <string>/Users/CHANGE_ME/Library/Logs/chatbotchat.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/CHANGE_ME/Library/Logs/chatbotchat.err.log</string>
</dict>
</plist>
